# linux-desktop/dot-local/sv/icloudpd/finish
#!/usr/bin/env sh
# wait 60 seconds before restart
echo "icloudpd/finish: waiting 60 seconds before restart"
sleep 60
