# linux-desktop/dot-local/services/cliphist/run
#!/usr/bin/env sh
if [[ -n $WAYLAND_DISPLAY ]]; then
  exec wl-paste --watch cliphist store
fi

echo "cliphist/run: not on wayland"
