# linux-desktop/services/niri-switch/run
#!/usr/bin/env sh
exec niri-switch-daemon
