# linux-desktop/dot-local/services/icloudpd/run
#!/usr/bin/env sh
exec icloudpd                      \
  --username moin@nilsriedemann.de \
  --directory ~/Pictures/Photos    \
  --watch-with-interval 3600
