# linux-desktop/services/waybar/run
#!/usr/bin/env sh
if [[ -n $WAYLAND_DISPLAY ]]; then
  exec waybar
  exit 0
fi

echo "waybar/run: not on wayland"
